# roles/globalprotect/templates/globalprotect.service.j2
[Unit]
Description={{ service_description }}

[Service]
Type=oneshot
ExecStart=/opt/globalprotect/reinstall.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target